# shithouse.init
#!/sbin/openrc-run
command="/usr/sbin/shithouse_tv"
command_args="-s /srv/shithouse"
supervisor="supervise-daemon"
